# discord.service
[Unit]
Description=DAO DAO Anti-Spam Discord Bot

[Install]
WantedBy=default.target

[Service]
ExecStart=/home/daodao/discord/run
Restart=always
